# 21-udev/src/90-usb-ether.rules
#-----------------------
# (c) H.Buchmann FHNW 2017
#adapt & copy /etc/udev/rules.d/90-usb-ether.rules
#----------------------- 

ACTION=="move", SUBSYSTEM=="usb"  \
RUN+="/home/buchmann/usb-eth.sh enable $env{INTERFACE} $env{ID_VENDOR_ID}"

ACTION=="remove", SUBSYSTEM=="usb"  \
RUN+="/home/buchmann/usb-eth.sh disable $env{INTERFACE} $env{ID_VENDOR_ID}"
=======
RUN+="/home/joel/Desktop/tinL1/tinL_own/21-udev/src/usb-eth.sh disable $env{INTERFACE} $env{ID_VENDOR_ID}"
=======

ACTION=="move", SUBSYSTEM=="net", ENV{ID_VENDOR_ID}=="0525" \
RUN+="/home/buchmann/usb-eth.sh enable $env{INTERFACE} $env{ID_VENDOR_ID}"

ACTION=="remove", SUBSYSTEM=="net", ENV{ID_VENDOR_ID}=="0525" \
RUN+="/home/buchmann/usb-eth.sh disable $env{INTERFACE} $env{ID_VENDOR_ID}"
>>>>>>> master
>>>>>>> Stashed changes
